# supervisord.conf
[program:dnsmasq]
command=dnsmasq -d

[program:docker-gen]
command=docker-gen -watch -notify "supervisorctl -c /etc/supervisor/supervisord.conf restart dnsmasq" /app/dnsmasq.tmpl /etc/dnsmasq.conf
